# roles/log/files/elasticsearch.conf
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
stopwaitsecs=3
stdout_logfile=/var/log/es_out.log
stderr_logfile=/var/log/es_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
