# roles/common/files/hekad.conf
description "hekad"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/hekad
